# lab-tanzu-platform-developer-experience/workshop/supervisor/cleanup-space.conf
[program:cleanup-space]
process_name=cleanup-space
command=/opt/workshop/supervisor/cleanup-space.sh
stdout_logfile=/proc/1/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true
